# crontab
0 4 * * * /bin/bash /usr/bin/image-sync.sh
